# rss-bot.service
[Unit]
Description=RSS Discord Bot Service
After=network.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user/rss-discord-bot
ExecStart=/usr/bin/python3 /home/ec2-user/rss-discord-bot/rss_discord_bot.py
Restart=on-failure
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target 