# rss-bot.timer
[Unit]
Description=Run RSS Bot on weekdays at 12:15 PM Denver time

[Timer]
OnCalendar=Mon..Fri 12:15:00 America/Denver
AccuracySec=1m
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target 