# rss-bot.timer
[Unit]
Description=Run RSS Bot daily at 12:15 PM Denver time

[Timer]
OnCalendar=*-*-* 12:15:00 America/Denver
AccuracySec=1m
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target 