# rss-bot.service
[Unit]
Description=RSS Bot Service
After=network.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user/rss-bot
ExecStart=/usr/bin/python3 /home/ec2-user/rss-bot/rss_monitor.py
Restart=on-failure
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target 